# collector/debian/gcm-sysmon-collector.cron.d
* * * * * root /usr/bin/sysmon-collect.py
* * * * * root /usr/bin/sysmon-checker.py