# collector/debian/gcm-sysmon-collector-core.cron.d
* * * * * root /usr/bin/sysmon-collect.py -f /etc/gcm-sysmon/sysmon.cfg
* * * * * root /usr/bin/sysmon-checker.py -f /etc/gcm-sysmon/sysmon.cfg

